# deploy/supervisor/siinv-exports.conf.example
; Copy this file to the server and replace every {{...}} placeholder.
; Run exactly one exports worker unless capacity testing justifies another process.
[program:siinv-exports]
process_name=%(program_name)s_%(process_num)02d
command={{PHP_BINARY}} {{APP_PATH}}/artisan queue:work database --queue=exports --sleep=3 --tries=2 --timeout=75 --memory=192 --max-jobs=100 --max-time=3600 --no-interaction
directory={{APP_PATH}}
user={{APP_USER}}
numprocs=1
autostart=true
autorestart=true
startsecs=5
startretries=3
stopasgroup=true
killasgroup=true
stopwaitsecs=105
redirect_stderr=true
stdout_logfile={{LOG_PATH}}/siinv-exports.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
